# config/supervisord.ini
[supervisord]
logfile = /var/log/supervisord.log

[program:application]
command = gunicorn --config /etc/opt/gunicorn.conf.py siteapp.wsgi
directory = /opt
stderr_logfile = /var/log/application-stderr.log
stdout_logfile = /var/log/application-stdout.log

[program:notificationemails]
command = ./manage.py send_notification_emails forever
directory = /opt
stderr_logfile = /var/log/notificationemails-stderr.log
stdout_logfile = /var/log/notificationemails-stdout.log
